# config/includes.chroot/etc/systemd/system/display-manager.service
[Unit]
Description=Masked - RBI boots to TTY
[Service]
ExecStart=/bin/true
